# gunicorn.socket
[Unit]
Description=<DESCRIPTION>

[Socket]
ListenStream=/run/<SOCK_NAME>

[Install]
WantedBy=sockets.target
