# uvicorn.socket
[Unit]
Description=<DESCRIPTION>

[Socket]
ListenStream=/run/<SOCK_NAME>.sock

[Install]
WantedBy=sockets.target
